# systemd/mgv_led_status.service
[Unit]
Description=MGV LED de estado RGB
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=Argos
WorkingDirectory=/home/Argos/MGV_Apolo_Epsilon_Base/src
ExecStart=/home/Argos/MGV_Apolo_Epsilon_Base/venv/bin/python /home/Argos/MGV_Apolo_Epsilon_Base/src/led_rgb_status.py
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
